# debian/igestis-ishare/etc/cron.d/ishare-fetchmail
*/5 * * * *	root	/usr/bin/ishare-fetchmail
